# theia/ide/admin/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log

[program:theia]
environment=HOME="/home/anubis",SHELL="/bin/bash",THEIA_DEFAULT_PLUGINS="local-dir:/opt/theia/plugins",DOCKER_CONFIG="/docker/"
command=/usr/local/bin/node /opt/theia/src-gen/backend/main.js /home/anubis/%(ENV_REPO_NAME)s --hostname=0.0.0.0 --port=5000
directory=/opt/theia
autorestart=true
user=anubis
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:dockerd]
command=/usr/bin/dockerd --storage-driver vfs
autorestart=true
user=root
stdout_logfile=/tmp/dockerd.log

[program:incluster-init]
directory=/home/anubis
user=anubis
command=/initialize-incluster.py
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:autosave-dump]
directory=/
user=anubis
command=/autosave-dump.sh
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
