# theia/ide/admin/supervisord.conf
[supervisord]
logfile=/supervisord.log

[program:hostname]
command=sh -c "echo anubis > /etc/hostname"
user=root
autorestart=false

[program:theia]
environment=HOME="/home/theia",SHELL="/bin/bash",THEIA_DEFAULT_PLUGINS="local-dir:/home/theia/plugins"
command=/usr/local/bin/node /home/theia/src-gen/backend/main.js /home/project --hostname=0.0.0.0 --port=5000
directory=/home/theia
autorestart=true
user=theia
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:dockerd]
command=/usr/bin/dockerd --storage-driver vfs
user=root
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

