# theia/ide/theia-golang/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log

[program:theia]
directory=/opt/theia
stopsignal=KILL
command=env -i REPO_NAME=%(ENV_REPO_NAME)s HOME="/home/anubis" THEIA_DEFAULT_PLUGINS="local-dir:/opt/theia/plugins/" SHELL="/bin/bash" PATH=%(ENV_PATH)s ANUBIS="1" ANUBIS_IDE="1" GO111MODULE=on /usr/local/bin/node /opt/theia/src-gen/backend/main.js /home/anubis/%(ENV_REPO_NAME)s --hostname=0.0.0.0 --port=5000
user=anubis
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:gopls]
directory=/opt/theia
command=env -i HOME=/home/anubis SHELL=/bin/bash PATH=%(ENV_PATH)s GO111MODULE=on ANUBIS=1 ANUBIS_IDE=1 /opt/go-tools/bin/gopls
user=anubis
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:autosave-dump]
directory=/
user=anubis
command=/autosave-dump.sh
autorestart=false
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0