# theia/ide/theia-base/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log

[program:theia]
directory=/opt/theia
stopsignal=KILL
command=env -i REPO_NAME="%(ENV_REPO_NAME)s" HOME="/home/anubis" THEIA_DEFAULT_PLUGINS="local-dir:/opt/theia/plugins/" SHELL="/bin/bash" PATH="/bin:/usr/bin:/usr/local/bin" ANUBIS="1" ANUBIS_IDE="1" /usr/local/bin/node /opt/theia/src-gen/backend/main.js /home/anubis/%(ENV_REPO_NAME)s --hostname=0.0.0.0 --port=5000
user=anubis
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:autosave-dump]
directory=/
user=anubis
command=/autosave-dump.sh
autorestart=false
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
